# usystem/etc/profile.d/env.entware
#!/usr/bin/env bash

# Loading EntWare if present
if [ -r /opt/etc/profile ] && bash -n /opt/etc/profile; then
    . /opt/etc/profile
fi
