# usystem/etc.defaults/rc.local
#!/bin/sh

[ -x /etc/rc.usystem ] && /etc/rc.usystem start

exit 0
